# backend/systemd/phillipsburg-radio-backend.timer
[Unit]
Description=Refresh Phillipsburg Radio backend JSON

[Timer]
OnBootSec=90s
OnUnitActiveSec=60s
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
